# apps/openstack/ansible/templates/cinder-volume-usage-audit.timer
[Unit]
Description=cinder-volume-usage-audit timer

[Timer]
OnBootSec=30min
OnCalendar=hourly
Persistent=True

[Install]
WantedBy=multi-user.target
